# launchd/com.maplion.claude-slack-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maplion.claude-slack-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/rozzum/git/maplion/maplion-ops-claude-slack-bridge/scripts/launchd-run.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/rozzum/git/maplion/maplion-ops-claude-slack-bridge</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/rozzum/Library/Logs/maplion-ops-claude-slack-bridge/out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/rozzum/Library/Logs/maplion-ops-claude-slack-bridge/err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
